# deploy/coding-server-user.service
[Unit]
Description=Coding Server (Flask/gunicorn) for classroom exercises
After=network.target

[Service]
Type=simple
# CHANGE ME: absolute path to where you cloned/deployed the app.
WorkingDirectory=/path/to/CodingServer
EnvironmentFile=/path/to/CodingServer/.env
# --access-logfile/--error-logfile are both required: gunicorn's access log
# (one line per request) is OFF by default. '-' sends both to stdout/stderr,
# which systemd captures into the journal automatically (see `journalctl` below).
# Bound to 127.0.0.1 (loopback only): nginx is the only thing that should
# talk to this directly -- see deploy/nginx-default.conf. No firewall rule
# is needed for this port as a result.
ExecStart=/path/to/CodingServer/.venv/bin/gunicorn --bind 127.0.0.1:8080 --access-logfile - --error-logfile - app:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target

# --- Install (no sudo required) ---
#   mkdir -p ~/.config/systemd/user
#   cp coding-server-user.service ~/.config/systemd/user/coding-server.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now coding-server
#
# --- Update after a deploy ---
#   systemctl --user restart coding-server
#
# --- Logs ---
#   journalctl --user -u coding-server -f
#
# --- IMPORTANT CAVEAT ---
# User-level systemd services normally stop when you log out, unless an
# administrator enables "lingering" for your account:
#   sudo loginctl enable-linger <your-username>
# That's a one-time favor to ask IT for. Without it, use option 3 in
# DEPLOYMENT.md (the cron/nohup fallback) instead, which survives logout
# and reboot with zero admin involvement.
