# deploy/coding-server.service
[Unit]
Description=Coding Server (Flask/gunicorn) for classroom exercises
After=network.target

[Service]
Type=simple
# CHANGE ME: the account this should run as on the dev server.
User=CHANGE_ME
# CHANGE ME: absolute path to where you cloned/deployed the app.
WorkingDirectory=/path/to/CodingServer
EnvironmentFile=/path/to/CodingServer/.env
# --access-logfile/--error-logfile are both required: gunicorn's access log
# (one line per request) is OFF by default. '-' sends both to stdout/stderr,
# which systemd captures into the journal automatically (see `journalctl` below).
# Bound to 127.0.0.1 (loopback only): nginx is the only thing that should
# talk to this directly -- see deploy/nginx-default.conf. No firewall rule
# is needed for this port as a result.
ExecStart=/path/to/CodingServer/.venv/bin/gunicorn --bind 127.0.0.1:8080 --access-logfile - --error-logfile - app:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

# --- Install (requires sudo/root) ---
#   sudo cp coding-server.service /etc/systemd/system/coding-server.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now coding-server
#
# --- Update after a deploy ---
#   sudo systemctl restart coding-server
#
# --- Logs ---
#   sudo journalctl -u coding-server -f
